# systemd/gatetorio-ble.service
[Unit]
Description=Gatetorio BLE Bridge Service
Documentation=file:///home/doowkcol/Gatetorio_Code/BLE_IMPLEMENTATION_GUIDE.md
After=bluetooth.target network.target
Wants=bluetooth.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -u /home/doowkcol/Gatetorio_Code/ble_server_bluezero.py
WorkingDirectory=/home/doowkcol/Gatetorio_Code
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=5
User=root

# Security settings (adjust as needed)
# Note: BLE requires root or bluetooth group membership
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW

[Install]
WantedBy=multi-user.target
